# packages/arch/ignite-bin/ignited.service
[Unit]
Description=ignited
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/ignited daemon --log-level debug

[Install]
WantedBy=multi-user.target
