# config/systemd/user/sway-session.target
[Unit]
Description=sway compositor session
Documentation=man:systemd.special(7)
BindsTo=graphical-session.target
# Run early session setup before the graphical session is declared up.
Wants=graphical-session-pre.target
After=graphical-session-pre.target
# Session daemons. Wants (not Requires) so a daemon crash never takes the
# session down; teardown is via each unit's own PartOf=graphical-session.target,
# not here (Wants is start-only). gnome-keyring first: it owns
# org.freedesktop.secrets and is unlocked at launch (it cannot be unlocked after
# start -- see gnome-keyring.service).
Wants=gnome-keyring.service
Wants=fcitx5.service
Wants=polkit-gnome.service
Wants=monitor-battery.timer
Wants=waybar.service
Wants=kanshi.service
